# scripts/com.openclaw.memory-maintenance-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.memory-maintenance-daily</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/juliocezar/Desenvolvimento/openclawdev/scripts/memory-maintenance.sh</string>
        <string>daily</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>59</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/memory-maintenance/daily.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/memory-maintenance/daily.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/juliocezar/.bun/bin</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/juliocezar/Desenvolvimento/openclawdev</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
